# rtl_433-mqtt/rootfs/etc/services.d/rtl_433-mqtt/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# rtl_433 to MQTT hass.io addon
# Runs rtl_433
# ==============================================================================

bashio::log.info "Starting rtl_433..."

#exec /usr/local/bin/rtl_433 -c /etc/rtl_433.conf -C si -M newmodel
exec /usr/local/bin/rtl_433 -d :10000002 -A
